# cfy_manager/components/postgresql_server/supervisord.conf
[program:cloudify-postgresql-server]
command=/usr/pgsql-9.5/bin/pg_ctl start -D /var/lib/pgsql/9.5/data/ -s -w -t 300
autostart=false
autorestart=false
startsecs=0
