# packaging/docker/starter.service
[Unit]
Wants=default.target
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "echo hello"
